# autostart/compton.desktop
[Desktop Entry]
Categories=Utility;
Comment=A X compositor
Exec=compton
GenericName=X compositor
Icon=compton
Name=compton
TryExec=compton
Type=Application
Version=1.0
