# autostart/imwheel.desktop
[Desktop Entry]
Exec=imwheel
Icon=
Name=imwheel
Path=
Terminal=False
Type=Application
